# .platform/files/etc/cron.hourly/cron.logrotate.elasticbeanstalk.cloudfront-api-key.conf
#!/bin/sh
test -x /usr/sbin/logrotate || exit 0
/usr/sbin/logrotate -f /etc/logrotate.elasticbeanstalk.hourly/logrotate.elasticbeanstalk.cloudfront-api-key.conf 